# apps/proxy/deploy/corrector-proxy.service
[Unit]
Description=Corrector Proxy API
After=network.target

[Service]
Type=simple
User=corrector
Group=corrector
WorkingDirectory=/opt/corrector/apps/proxy
EnvironmentFile=/etc/corrector/proxy.env
ExecStart=/usr/bin/node /opt/corrector/apps/proxy/src/server.js
Restart=always
RestartSec=2
LimitNOFILE=65535
TimeoutStopSec=20
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
